# ansible/roles/etcd/templates/etcd.service.j2
[Unit]
Description=etcd key-value store
Documentation=https://github.com/etcd-io/etcd
After=network.target
Wants=network-online.target

[Service]
Type=notify
User={{ etcd_user }}
Group={{ etcd_group }}
ExecStart={{ etcd_bin_dir }}/etcd \
  --name={{ inventory_hostname }} \
  --data-dir={{ etcd_data_dir }} \
  --listen-client-urls={{ etcd_listen_client_urls }} \
  --advertise-client-urls={{ etcd_advertise_client_urls }} \
  --listen-peer-urls={{ etcd_listen_peer_urls }} \
  --initial-advertise-peer-urls={{ etcd_initial_advertise_peer_urls }} \
  --initial-cluster={{ etcd_initial_cluster }} \
  --initial-cluster-token={{ etcd_initial_cluster_token }} \
  --initial-cluster-state={{ etcd_initial_cluster_state }} \
  --client-cert-auth=true \
  --trusted-ca-file={{ etcd_ca_file }} \
  --cert-file={{ etcd_cert_dir }}/etcd-{{ inventory_hostname }}.pem \
  --key-file={{ etcd_cert_dir }}/etcd-{{ inventory_hostname }}-key.pem \
  --peer-client-cert-auth=true \
  --peer-trusted-ca-file={{ etcd_ca_file }} \
  --peer-cert-file={{ etcd_cert_dir }}/etcd-{{ inventory_hostname }}.pem \
  --peer-key-file={{ etcd_cert_dir }}/etcd-{{ inventory_hostname }}-key.pem \
  --log-level=info \
  --logger=zap \
  --log-outputs=systemd/journal
Restart=on-failure
RestartSec=5
LimitNOFILE=40000

[Install]
WantedBy=multi-user.target